# ops/systemd/ops-validate.timer
[Unit]
Description=Weekly backup validation drill

[Timer]
OnCalendar=Sat 04:10
Persistent=true
Unit=ops-validate.service

[Install]
WantedBy=timers.target
